# debian_stuff/crontab
- run every minute
*/1 * * * * /path/to/exec

- run every hour
* */1 * * * /path/to/exec
